# playbook/config/vault/vault.service
[Unit]
Description=Vault service

[Service]
ExecStart=/usr/bin/vault server -config=/etc/opt/vault/config.hcl

[Install]
WantedBy=multi-user.target